# com.immo-search-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.immo-search-weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/wallaceespindola/git/immo-search/run.sh</string>
        <string>--week</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/wallaceespindola/git/immo-search</string>

    <!-- Every Saturday at 09:00 (Weekday 6 = Saturday in launchd) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>6</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/wallaceespindola/git/immo-search/logs/weekly-launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/wallaceespindola/git/immo-search/logs/weekly-launchd.error.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
